# arco-virt/repo/script/inittab
::sysinit:/etc/rc.d/rc.sysinit
::respawn:/sbin/getty 115200 ttyAMA0
